# port/files/gwman_agent.in
#!/bin/sh
#
# PROVIDE: gwman_agent
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# gwman_agent_enable  (bool):  Set to NO by default.
#               Set it to YES to enable gwman_agent.
#

. /etc/rc.subr

name="gwman_agent"
rcvar="${name}_enable"
command="/usr/local/bin/gwman_agent"
command_interpreter="/usr/local/bin/python2.7"

load_rc_config $name

: ${gwman_agent_enable="NO"}
: ${gwman_agent_config="/usr/local/etc/gwman/agent.conf"}
: ${gwman_agent_pidfile="/var/run/gwman_agent.pid"}

pidfile="${gwman_agent_pidfile}"

start_cmd="${command} start -c ${gwman_agent_config}"
stop_cmd="${command} stop -c ${gwman_agent_config}"
restart_cmd="${command} restart -c ${gwman_agent_config}"

run_rc_command $1
