# port/files/gwman_server.in
#!/bin/sh
#
# PROVIDE: gwman_server
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# gwman_server_enable  (bool):  Set to NO by default.
#               Set it to YES to enable gwman_server.
#

. /etc/rc.subr

name="gwman_server"
rcvar="${name}_enable"
command="/usr/local/bin/gwman_server"
command_interpreter="/usr/local/bin/python2.7"

load_rc_config $name

: ${gwman_server_enable="NO"}
: ${gwman_server_config="/usr/local/etc/gwman/server.conf"}
: ${gwman_server_pidfile="/var/run/gwman_server.pid"}

pidfile="${gwman_server_pidfile}"

start_cmd="${command} start -c ${gwman_server_config}"
stop_cmd="${command} stop -c ${gwman_server_config}"
restart_cmd="${command} restart -c ${gwman_server_config}"

run_rc_command $1
